# ops/cicf/cicf-bridge.service
[Unit]
Description=CICF Transparent Bridge (br0)
After=network-pre.target
Wants=network-pre.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cicf-bridge.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

